# ops/eu.verwey.reck-mount.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>eu.verwey.reck-mount</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/reck-mount-watchdog</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <!-- RECK_STATION_ROOT is the path on the station where managed
         projects live (e.g. /Users/reck-connect/projects for the
         upstream macOS station, /home/<user>/projects for a Linux
         station). install-satellite.sh substitutes __RECK_STATION_ROOT__
         with the operator-supplied value at install time; launchd does
         not expand env vars in this dict, so the value is baked in
         literally. The watchdog (reck-mount-watchdog.sh) reads it via
         shell env on every tick and fails-fast if absent. -->
    <key>RECK_STATION_ROOT</key>
    <string>__RECK_STATION_ROOT__</string>
  </dict>
  <!-- Logs live under ~/Library/Logs/reck-stationd/ — per-user (parent
       home is 0700) instead of the world-readable /tmp, and they
       survive reboot. install-satellite.sh substitutes __HOME__ for
       the actual $HOME when it copies this template into place, since
       launchd does NOT expand environment variables in path keys. -->
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/reck-stationd/mount.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/reck-stationd/mount.log</string>
</dict>
</plist>
